# scripts/rio-coffee.service
[Unit]
Description=Rio Coffee Next.js app
After=network.target

[Service]
Type=simple
User=www-data
# The installer script will replace __WORKDIR__ with the repository absolute path
WorkingDirectory=__WORKDIR__
Environment=NODE_ENV=production
# Use npx from PATH; installer will ensure WorkingDirectory is correct
ExecStart=/usr/bin/env npx next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rio-coffee

[Install]
WantedBy=multi-user.target
